# BHOI_/2023/bonbon/IN-OUT/crontab
0 14 * * * /home/cmnikolas/backup_foodfoca.sh
